# tovuq.conf
[program:tovuqbot]

directory=/home/ubuntu/qora_tovuq
command=python3 app.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu",USER="ubuntu"

stderr_logfile=/home/ubuntu/qora_tovuq/logfile_err.log
stdout_logfile=/home/ubuntu/qora_tovuq/logfile.log